# deploy/mcp-e2e/mcp-e2e.service.example
# 示例：在测试机上对已部署的 trade-mcp 做 MCP 协议探活
# 依赖：/opt/AIFutureTrade/deploy/mcp-e2e 下已 pip install -r requirements-e2e.txt
#
# /etc/mcp-e2e.env 示例：
#   TRADE_MCP_BASE_URL=http://127.0.0.1:8099

[Unit]
Description=AIFutureTrade trade-mcp MCP (SSE) end-to-end probe
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/mcp-e2e.env
ExecStart=/usr/bin/python3 /opt/AIFutureTrade/deploy/mcp-e2e/mcp_e2e.py
WorkingDirectory=/opt/AIFutureTrade/deploy/mcp-e2e

[Install]
WantedBy=multi-user.target
